# build_files/services/games.service
[Unit]
Description="Decrypts and mount games disks"

[Service]
ExecStart=/usr/lib/systemd/system/games/games.sh
