# build_files/services/games.service
[Unit]
Description="Decrypts and mount games disks"

[Service]
ExecStart=/etc/systemd/system/games.sh
